# .config/systemd/user/start.service
[Unit]
Description=Start

[Service]
Type=oneshot
ExecStart=/bin/bash /home/qam1s/.start.sh

[Install]
WantedBy=default.target
